# files/systemd/user/i3.target
[Unit]
Description=i3 running
BindsTo=xsession.target
