# files/systemd/user/pyi3wsd.service
[Unit]
Description=i3 Workspace Daemon

[Service]
ExecStart=%h/.local/bin/pyi3wsd
Restart=always

[Install]
WantedBy=xsession.target
