# files/systemd/user/pyi3wsd.service
[Unit]
Description=i3 Workspace Daemon
PartOf=i3.target

[Service]
ExecStart=%h/.local/bin/pyi3wsd
Restart=always

[Install]
WantedBy=i3.target
